# package/multimedia/pulseaudio/init_scripts/pulseaudio.init
#!/bin/sh /etc/rc.common

START=99
STOP=65

USE_PROCD=1
PROG=/usr/bin/pulseaudio

start_service() {
	procd_open_instance
	procd_set_param command $PROG --system --disallow-exit --disallow-module-loading --disable-shm --exit-idle-time=-1 --realtime=false
	procd_close_instance
}

stop_service() {
	[ -e /var/run/pulse/pid ] && kill $(cat /var/run/pulse/pid)
}

