# package/testtools/tinatest/files/tinatest.init
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1
CMD="tinatest"
DEPEND="sysntpd"

start_service() {
	procd_open_instance
	procd_set_param command ${CMD}
	procd_set_param respawn
	procd_close_instance
}
