# package/allwinner/boot-play/play.init
#!/bin/sh /etc/rc.common

START=1
STOP=98
DEPEND=done
PROG=boot-play

start() {
	$PROG boot &
}

stop() {
	$PROG shutdown
}

