# package/busybox-init-base-files/busybox-init-base-files_ramfs/etc/inittab
::sysinit:/etc/init.d/rcS boot
::respawn:-/bin/sh
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
