# package/ros/ros-base/ros_comm/files/rossetupenv.sh
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2019 OpenWrt.org

export CMAKE_PREFIX_PATH=/usr
export PKG_CONFIG_PATH=/usr/lib/pkgconfig
export ROS_ROOT=/usr/share/ros
export ROS_PACKAGE_PATH=/usr/share
export ROS_LOG_DIR=/mnt/UDISK
export ROS_HOSTNAME=localhost
export ROS_MASTER_URI=http://localhost:11311
export LD_LIBRARY_PATH=/usr/lib
export ROS_ETC_DIR=/usr/etc/ros
export ROS_DISTRO=kinetic
export ROS_VERSION=1
export ROSLISP_PACKAGE_DIRECTORIES=/usr/share/common-lisp
